# omniagent.service
[Unit]
Description=OmniAgent v8.1 — Autonomous AI Agent Framework
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=kraken
WorkingDirectory=/home/kraken/Devel/OmniAgent
ExecStart=/usr/bin/python3 /home/kraken/Devel/OmniAgent/omni_agent.py
Restart=on-failure
RestartSec=5
Environment=PATH=/home/kraken/Devel/OmniAgent/.venv/bin:/usr/local/bin:/usr/bin
Environment=OLLAMA_NUM_CTX=32768

[Install]
WantedBy=multi-user.target
